# Resouces/etc/profile.d/autostart.sh
# Auto start Hyprland at login
Hyprland

# Auto start kitty at login
sleep 3
kitty

