# Resouces/etc/profile.d/autostart.sh
# Auto start Hyprland at login
Hyprland

ibus-daemon -rxRd
