# Resouces/etc/profile.d/start-ibus.sh
# Auto start ibus at login
ibus-daemon -rxRd
