# Resouces/etc/profile.d/start-hyprland.sh
# Auto start Hyprland at login
Hyprland

